# deploy/supervisor-worker.conf
; ─────────────────────────────────────────────────────────────────────────────
; M2 Guardian — Queue Worker (Supervisor)
; Copiar para: /etc/supervisor/conf.d/m2guardian-worker.conf
; Depois:
;   supervisorctl reread
;   supervisorctl update
;   supervisorctl start m2guardian-worker:*
; ─────────────────────────────────────────────────────────────────────────────

[program:m2guardian-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/m2guardian/artisan queue:work database --sleep=3 --tries=3 --max-time=3600 --timeout=90
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=m2guardian
numprocs=2
redirect_stderr=true
stdout_logfile=/var/log/m2guardian-worker.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stopwaitsecs=120
